# meta-bsp-tv/recipes-bsp/powervr/omap4-sgx-modules/rc.pvr
#!/bin/sh
### BEGIN INIT INFO
# Provides:          pvr-init
# Required-Start:    checkroot
# Required-Stop:
# Default-Start:     S
# Default-Stop:
# Short-Description: PVR SGX driver init.
### END INIT INFO

[ -e /sbin/modprobe ] || exit 1

LOAD_MODULE=/sbin/modprobe

[ "$VERBOSE" != no ] && echo -n "Loading PVR modules: "

$LOAD_MODULE omapdrm_pvr

[ -e /usr/bin/pvrsrvinit ] && /usr/bin/pvrsrvinit
